# modules/nixos/workflows/mosey-branch/config/desktop-session/core-service.target.in
[Unit]
Description=@description@
Documentation=man:gnome-session(1) man:systemd.special(7)
CollectMode=inactive-or-failed

# Pull in the service.
Wants=@app_id@.@component@.service

# Require GNOME session and do some reordering.
PartOf=gnome-session-initialized.target
After=gnome-session-initialized.target
Before=gnome-session.target
